# engines/data/user/projets/tco-core/gateway/gateway.service
[Unit]
Description=TCO Lab Gateway - Acquisition Arduino vers InfluxDB
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=tco
Group=tco
# Le dossier où se trouve ton script sur le Pi
WorkingDirectory=/home/tco/tco-core/gateway
ExecStart=/home/tco/tco-core/gateway/venv/bin/python3 /home/tco/tco-core/gateway/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tco-gateway

# Sécurité
NoNewPrivileges=true
PrivateTmp=true

# Permissions pour accéder au port USB
SupplementaryGroups=dialout

[Install]
WantedBy=multi-user.target
